# service/node/run
#!/bin/sh

exec /usr/bin/node ~/build/app.js
